# envs/aws_linux/systemd/nologin_shutdown.timer
[Unit]
Description=Periodically check who is logged in; shutdown if no one

[Timer]
OnCalendar=*-*-* *:0/10

[Install]
WantedBy=default.target
