# systemd/nologin_shutdown.service
[Unit]
Description=Shutdown this server if no one is logged in

[Service]
ExecStart=/home/njormrod/dev/devenv/bin/shutdown_if_no_login.py
